# Service/waterbot.service
[Unit]
Description=Waterbot daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/Waterbot -daemon

[Install]
WantedBy=multi-user.target
